# deployment/celery-beat.service
[Unit]
Description=Celery Beat Scheduler Service for RS Systems Notification System
After=network.target redis.service postgresql.service celery-worker.service
Requires=celery-worker.service

[Service]
Type=simple
User=django
Group=django
WorkingDirectory=/var/www/rs_systems

# Environment variables
Environment="PATH=/var/www/rs_systems/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=rs_systems.settings_aws"
Environment="ENVIRONMENT=production"

# Load additional environment variables from file (optional)
# EnvironmentFile=/var/www/rs_systems/.env

# Celery beat command
ExecStart=/var/www/rs_systems/venv/bin/celery -A rs_systems beat \
    --loglevel=info \
    --pidfile=/var/run/celery/beat.pid \
    --logfile=/var/log/celery/beat.log \
    --schedule=/var/run/celery/celerybeat-schedule

# Graceful shutdown
ExecStop=/bin/kill -s TERM $MAINPID

# Runtime directory
RuntimeDirectory=celery
RuntimeDirectoryMode=0755

# Logging directory
LogsDirectory=celery
LogsDirectoryMode=0755

# Process management
PIDFile=/var/run/celery/beat.pid
Restart=always
RestartSec=10s

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
