# deployment/celery-worker.service
[Unit]
Description=Celery Worker Service for RS Systems Notification System
After=network.target redis.service postgresql.service

[Service]
Type=forking
User=django
Group=django
WorkingDirectory=/var/www/rs_systems

# Environment variables
Environment="PATH=/var/www/rs_systems/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE=rs_systems.settings_aws"
Environment="ENVIRONMENT=production"

# Load additional environment variables from file (optional)
# EnvironmentFile=/var/www/rs_systems/.env

# Celery worker command
ExecStart=/var/www/rs_systems/venv/bin/celery -A rs_systems worker \
    --loglevel=info \
    --concurrency=8 \
    --queues=celery,notifications,maintenance \
    --pidfile=/var/run/celery/worker.pid \
    --logfile=/var/log/celery/worker.log \
    --detach

# Graceful shutdown
ExecStop=/bin/kill -s TERM $MAINPID

# Runtime directory
RuntimeDirectory=celery
RuntimeDirectoryMode=0755

# Logging directory
LogsDirectory=celery
LogsDirectoryMode=0755

# Process management
PIDFile=/var/run/celery/worker.pid
Restart=always
RestartSec=10s

# Resource limits (optional but recommended)
LimitNOFILE=65536
MemoryLimit=2G
CPUQuota=400%

[Install]
WantedBy=multi-user.target
